# packages/php-fpm-exporter/php-fpm-exporter.service
[Unit]
Description=A prometheus exporter for PHP-FPM.
Documentation=https://github.com/bakins/php-fpm-exporter
After=network.target

[Service]
User=monitoring
EnvironmentFile=-/etc/default/php-fpm-exporter
ExecStart=/usr/bin/php-fpm-exporter server $PHP_FPM_EXPORTER_OPTS
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
