# packages/elasticsearch-hprof/elasticsearch-crash-handler.path
[Path]
PathExistsGlob=/var/lib/elasticsearch/*.hprof
Unit=elasticsearch-crash-handler.service
PartOf=elasticsearch.service

[Install]
WantedBy=multi-user.target
