# packages/elasticsearch-hprof/elasticsearch-crash-handler.service
[Unit]
Description=Process Elasticsearch OOM heap dump

[Service]
Type=oneshot
WorkingDirectory=/var/lib/elasticsearch
ExecStart=/bin/bash -c 'for dump in timeout 15m /var/lib/elasticsearch/*.hprof;do /usr/bin/elasticsearch-hprof inflight_queries --save $dump; rm $dump;done'
