# packages/varnish-exporter/varnish-exporter.service
[Unit]
Description=Varnish exporter for Prometheus
Documentation=https://github.com/jonnenauha/prometheus_varnish_exporter
After=network.target

[Service]
User=monitoring
Group=varnish
EnvironmentFile=-/etc/default/varnish-exporter
ExecStart=/usr/bin/prometheus_varnish_exporter $VARNISH_EXPORTER_OPTS
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
